# walkie-openclaw-adapter/references/walkie-adapter.service.example
[Unit]
Description=Walkie OpenClaw Adapter
After=network.target

[Service]
Type=simple
# TODO: Adjust paths below for your environment
WorkingDirectory=YOUR_WORKSPACE_PATH
ExecStart=/path/to/node skills/walkie-openclaw-adapter/scripts/adapter.js skills/walkie-openclaw-adapter/references/config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# Example: Environment="PATH=/home/user/.nvm/versions/node/v22.22.0/bin:$PATH"
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
